# files/etc/profile.d/20-exports.sh
# shellcheck shell=sh
export EDITOR=nvim
export DELTA_PAGER="less -R"
